# KASL/lista9/mystat/mystat-graph.service
[Unit]
Description=Generate daily CPU usage graph

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mystat-graph.sh
